# .myScripts/crontab
0 0-8 * * * env DISPLAY=:0 sh /home/naash/.myScripts/blank_screen.sh > /dev/null 2>&1
